# supervisord.conf
[supervisord]
nodaemon=true

[program:dockerd]
command=/usr/local/bin/dockerd

[program:executor]
command=/app/executor